# .config/systemd/user/avfs.service
[Unit]
Description=AVFS Daemon

[Service]
Type=forking
ExecStart=/usr/bin/avfsd %h/.avfs

[Install]
WantedBy=default.target
