# .config/systemd/user/mako.service
[Unit]
Description=Notification server

[Service]
Type=simple
ExecStart=/usr/bin/mako

[Install]
WantedBy=default.target
