# .config/systemd/backup.timer
[Unit]
Description= Backup Timer

[Timer]
Unit=backup.service
OnBootSec=120
OnCalendar=hourly

[Install]
WantedBy=timers.target
